# agent/fpr-amd-agent.service
[Unit]
Description=Floating-Point Ranger — AMD GPU state reporter (read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/fpr-amd-agent/fpr_amd_agent.py
Restart=on-failure
RestartSec=5
# Runs as root INSIDE the unprivileged container: reading another user's
# /proc/<pid>/fdinfo/<fd> requires PTRACE_MODE_READ, and the GPU clients typically run as
# other users. Container-root is namespaced (an unprivileged uid on the host).
User=root
# The agent only ever READS /proc and /sys — lock the rest down.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryMax=128M

[Install]
WantedBy=multi-user.target
